# mail/vbsfilter/pkg-plist
@comment $FreeBSD$
libexec/vbsfilter
etc/rc.d/vbsfilter.sh.sample
