# net/isc-dhcp30-relay/pkg-plist
@comment $FreeBSD$
etc/rc.d/isc-dhcrelay.sh
sbin/dhcrelay
